# files/rc.local
#!/bin/sh

# Enter custom commands here. By default, they will be executed on startup in parallel with the other services.

# AC MODE
if [ -x /bin/ac-mode ]; then
    ( sleep 20; setsid nohup /bin/ac-mode >> /var/log/ac-mode.log 2>&1 ) &
fi

